# deploy/systemd/crypto-bot-runtime.service
# Copy this file to /etc/systemd/system/crypto-bot-runtime.service
# Adjust User, WorkingDirectory, EnvironmentFile, and ExecStart if your deployment path differs.

[Unit]
Description=Crypto Bot V1 Trading Runtime
Wants=network-online.target
After=network-online.target
# Prefer crypto-bot-tmux.service on VPS when you need Rich CLI + attach; do not enable both.
Conflicts=crypto-bot-tmux.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/Crypto_Sniper
EnvironmentFile=-/root/Crypto_Sniper/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8
Environment=PYTHONUTF8=1
Environment=LANG=en_US.UTF-8
Environment=LC_ALL=en_US.UTF-8
ExecStart=/root/Crypto_Sniper/.venv/bin/python /root/Crypto_Sniper/main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=30
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
