# deploy/systemd/crypto-bot-runtime.service
# Copy this file to /etc/systemd/system/crypto-bot-runtime.service
# Adjust User, WorkingDirectory, EnvironmentFile, and ExecStart if your deployment path differs.

[Unit]
Description=Crypto Bot V1 Trading Runtime
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/Crypto_Sniper
EnvironmentFile=-/root/Crypto_Sniper/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/Crypto_Sniper/.venv/bin/python /root/Crypto_Sniper/main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=30
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
