# deploy/systemd/crypto-bot-tmux.service
# Copy this file to /etc/systemd/system/crypto-bot-tmux.service
# Adjust User, Group, and environment paths if your deployment path differs.

[Unit]
Description=Crypto Bot V1 tmux Runtime
Wants=network-online.target
After=network-online.target
# Starting this unit stops direct-Python bots so only one main.py runs (inside tmux).
Conflicts=crypto-bot-runtime.service crypto-sniper.service

[Service]
Type=oneshot
User=root
Group=root
Environment=TMUX_SESSION_NAME=crypto
Environment=CRYPTO_BOT_ROOT=/root/Crypto_Sniper
Environment=CRYPTO_BOT_PYTHON=/root/Crypto_Sniper/.venv/bin/python
Environment=CRYPTO_BOT_ENTRY=/root/Crypto_Sniper/main.py
ExecStart=/root/Crypto_Sniper/deploy/systemd/crypto-bot-tmux.sh start
ExecStop=/root/Crypto_Sniper/deploy/systemd/crypto-bot-tmux.sh stop
ExecReload=/root/Crypto_Sniper/deploy/systemd/crypto-bot-tmux.sh restart
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target